# minecraft-stop@.service
# Service to stop the servers with the timers

[Unit]
Description=Stop Minecraft Server %i

[Service]
Type=oneshot

ExecStart=systemctl stop minecraft@%i
